# deploy/xpresspro.service
[Unit]
Description=XpressPro FX API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust the User and WorkingDirectory for your VPS environment
User=www-data
WorkingDirectory=/var/www/Rebranded-xpfx
EnvironmentFile=-/etc/xpressfx.env
Environment="NODE_ENV=production"
Environment="PORT=3000"

# The start script will install dev deps for the build step when necessary
ExecStart=/usr/bin/env bash /var/www/Rebranded-xpfx/start.sh
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=10

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/Rebranded-xpfx /var/log/xpresspro
LimitNOFILE=65536
LimitNPROC=32768
MemoryMax=2G

[Install]
WantedBy=multi-user.target
